# com.simah.tablestatus.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd agent: runs the table-status check cycle once a day at 03:00.
  Install (test first!):
    cp "com.simah.tablestatus.plist" ~/Library/LaunchAgents/
    launchctl load  ~/Library/LaunchAgents/com.simah.tablestatus.plist   # enable
    launchctl start com.simah.tablestatus                                # run now (test)
    launchctl unload ~/Library/LaunchAgents/com.simah.tablestatus.plist  # disable
  Requires: the capture Chrome open with the Teams 'table status' tab active,
  the VDI mail-watcher alive, and the Mac awake at 03:00 (else it runs at next wake).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.simah.tablestatus</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/alexandrgruzdev/Downloads/SIMAH/claude playgraund/table-status-checker/run_cycle.py</string>
        <string>--run</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/alexandrgruzdev/Downloads/SIMAH/claude playgraund/table-status-checker</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/alexandrgruzdev/Downloads/SIMAH/claude playgraund/table-status-checker/logs/cycle.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/alexandrgruzdev/Downloads/SIMAH/claude playgraund/table-status-checker/logs/cycle.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/bin:/bin:/usr/local/bin</string>
    </dict>
</dict>
</plist>
